# nvchecker/run-nvchecker.cron
#!/bin/bash

PYTHON="python3.9"
CHECKOUT="$HOME/archlinux-proaudio"

export PATH="$HOME/.local/bin:/usr/bin:/bin"

cd "$CHECKOUT"
git stash
git pull -q
git stash pop
nvchecker -l warning
$PYTHON nvchecker/nvchecker-notify-matrixchat.py # --dry-run
